# debian/update-motd.d/00-header
#!/bin/sh

uname -a
printf "%s %s (%s)\n" "$(lsb_release -s -i)" "$(lsb_release -s -r)" "$(lsb_release -s -c)"
